# service/uploadfilter.service
[Service]
ExecStart=/bin/SFTP-Filter/SysMon.sh
